# pidrn.service
[Unit]
Description=PI Drone Service
[Service]
ExecStart=/usr/bin/pi_drone
[Install]
WantedBy=multi-user.target
